# tools/ci/config/supervisord/supervisord.ini
[supervisord]
user=root
nodaemon=true
stdout_logfile=/proc/self/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0

[program:php-fpm]
command = /usr/sbin/php-fpm7 -F --nodaemonize
autostart=true
autorestart=true
stdout_logfile=/proc/self/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0
exitcodes=0

[program:nginx-server]
command = /usr/sbin/nginx
priority=900
stdout_logfile=/proc/self/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0
username=nginx
autorestart=true
